# infrastructure/ansible/templates/plantocode-boot-sync.service.j2
[Unit]
Description=Sync plantocode nginx upstream with active deployment on boot
After=network.target
Before=nginx.service plantocode-blue.service plantocode-green.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/plantocode-boot-sync
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
