# infrastructure/ansible/templates/plantocode.service.j2
[Unit]
Description={{ app_name }} Rust Service
After=network.target postgresql@17-main.service
Wants=postgresql@17-main.service

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}

# Binary path
ExecStart={{ app_dir }}/bin/server

# Graceful reload using symlink switch
ExecReload=/bin/kill -s HUP $MAINPID

# Restart policy
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ app_dir }}/logs /var/log/{{ app_name }}

# Resource limits
LimitNOFILE=65535

# Environment
EnvironmentFile={{ app_dir }}/config/app.env

# Health check disabled - using Type=simple
# WatchdogSec=10s

[Install]
WantedBy=multi-user.target
