# infrastructure/ansible/templates/vibe-manager-blue.service.j2
[Unit]
Description=vibe-manager Rust Service (Blue)
After=network.target postgresql@17-main.service
Wants=postgresql@17-main.service

[Service]
Type=simple
User=vibe-manager
Group=vibe-manager
WorkingDirectory=/opt/vibe-manager

# Binary path
ExecStart=/opt/vibe-manager/bin/server

# Environment override for blue instance
Environment="SERVER_PORT=8080"
Environment="DEPLOYMENT_COLOR=blue"
Environment="VIBE_DEPLOYMENT_TOKEN={{ vault_vibe_deployment_token }}"

# Graceful reload using symlink switch
ExecReload=/bin/kill -s HUP $MAINPID

# Restart policy
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vibe-manager/logs /var/log/vibe-manager

# Resource limits
LimitNOFILE=65535

# Environment file (contains all other env vars)
EnvironmentFile=/opt/vibe-manager/config/app.env

# Graceful shutdown timeout matching application
TimeoutStopSec=130
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes

[Install]
WantedBy=multi-user.target